# deploy/ppmap.service
[Unit]
Description=PPMap Node.js server
After=network.target mongod.service

[Service]
User=ppmap
WorkingDirectory=/srv/ppmap/PPMap/server
ExecStart=/usr/bin/node index.js
Restart=on-failure
RestartSec=5
EnvironmentFile=/srv/ppmap/PPMap/server/.env

[Install]
WantedBy=multi-user.target
