# deploy/ppmap-cert-renew.timer
[Unit]
Description=Renew PPMap TLS cert weekly

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
